# services/phrase-picker-swift/run
#!/bin/bash

set -x

swift run App serve --port 11111
